# src/maven/com.alis.maven.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alis.maven</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/alis_mini/.pyenv/versions/3.11.11/bin/python3</string>
        <string>/Users/alis_mini/.openclaw/workspace/scripts/maven/maven.py</string>
        <string>--once</string>
    </array>
    
    <key>StartInterval</key>
    <integer>300</integer>
    
    <key>WorkingDirectory</key>
    <string>/Users/alis_mini/.openclaw/workspace/scripts/maven</string>
    
    <key>StandardOutPath</key>
    <string>/tmp/openclaw/maven.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/maven.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/alis_mini/.pyenv/versions/3.11.11/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/alis_mini</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <!-- 설치:
    cp com.alis.maven.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.alis.maven.plist
    
    제거:
    launchctl unload ~/Library/LaunchAgents/com.alis.maven.plist
    rm ~/Library/LaunchAgents/com.alis.maven.plist
    
    수동 실행:
    launchctl start com.alis.maven
    
    로그 확인:
    tail -f /tmp/openclaw/maven.log
    -->
</dict>
</plist>
